# .supervisor/supervisord.conf
[supervisord]
logfile = /home/sirex/.supervisor/supervisord.log ; (main log file;default $CWD/supervisord.log)
pidfile = /home/sirex/.supervisor/supervisord.pid ; (supervisord pidfile;default supervisord.pid)
childlogdir = /home/sirex/.supervisor/logs        ; ('AUTO' child log dir, default $TEMP)

[unix_http_server]
file = /home/sirex/.supervisor/supervisor.sock    ; (the path to the socket file)
chmod = 0700                                      ; sockef file mode (default 0700)

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///home/sirex/.supervisor/supervisor.sock ; use a unix:// URL  for a unix socket

[program:bootstrap]
command = /home/sirex/bin/supervisor-bootstrap
autostart = true

[program:offlineimap]
command = offlineimap -u quiet
stopsignal = USR2
autostart = false

[program:encfs]
command = /home/sirex/bin/encfs-daemon
autostart = false
